# templates/sentry-worker.service.erb
[Unit]
Description=Sentry Worker
After=network.target

[Service]
User=<%= @user %>
Group=<%= @group %>
Type=simple
PIDFile=<%= @venv %>/sentry-worker.pid
Environment=VIRTUAL_ENV="<%= @venv %>"
Environment=PATH="$VIRTUAL_ENV/bin:$PATH"
WorkingDirectory=<%= @venv %>
ExecStart=<%= @venv %>/bin/sentry --config=<%= @venv %>/sentry.conf.py celery worker -B -f /var/log/sentry/sentry-celery.log --pidfile=<%= @venv %>/sentry-worker.pid
ExecStop=/bin/kill -KILL $MAINPID
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
